# lsb/xllmnrd.in
#!/bin/sh
# Init script for xllmnrd
# Copyright (C) 2013 Kaz Nishimura
#
# Copying and distribution of this file, with or without modification, are
# permitted in any medium without royalty provided the copyright notice and
# this notice are preserved.  This file is offered as-is, without any
# warranty.

### BEGIN INIT INFO
# Provides: xllmnrd llmnrd
# Required-Start: $local_fs $network $syslog
# Required-Stop: $local_fs $network $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start or stop the IPv6 LLMNR responder
# Description:
### END INIT INFO

. /lib/lsb/init-functions

log_warning_msg $0 is a placeholder script

exit 0
